# rasa/files/rasa-api.service
[Unit]
Description=Start the rasa web server API

[Service]
ExecStart=/opt/rasa/init_rasa_api.sh
Restart=always

[Install]
WantedBy=multi-user.target